# HW-lesson-09/dop_files/logs-nginx.cron
#!/bin/bash
/root/scripts/logs-nginx.sh
